# deploy/php/supervisor.conf
[supervisord]
nodaemon=true
user=root

# laravelの開発用サーバーを起動
[program:laravel-server]
command=php /var/www/html/artisan serve --host=0.0.0.0 --port=80 --env=.env
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-server.log

# 非同期タスク用キューのワーカーを起動
[program:queue-worker]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3
autostart=true
autorestart=true
user=root
